# roles/pinger/templates/pinger.service
[Unit]
After=network.target

[Service]
User=pinger
ExecStart=/usr/local/bin/pinger

[Install]
WantedBy=multi-user.target
